# src/kactivitymanagerd.service.in
[Unit]
Description=KActivityManager

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/kactivitymanagerd
KillMode=process
BusName=org.kde.ActivityManager
